# scripts/kvn-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kvn-web.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kvn-web</string>
        <string>--no-browser</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/kvn-web.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/kvn-web.log</string>
</dict>
</plist>
